# etc/init/gamblacoin.conf
description "gamblacoin node.js server"
author      "gamblacoin"

start on startup
stop on shutdown

script
	export HOME="/home/gamblacoin/"
	echo $$ > /var/run/gamblacoin.pid
	exec su -c "export NODE_ENV=prod-local && cd /home/gamblacoin/gamblacoinbuild/latest && /usr/local/bin/node app.js >> /home/gamblacoin/gamblacoin.sys.log 2>&1" - gamblacoin
end script

pre-start script
# Date format same as (new Date()).toISOString() for consistency
	echo "[`date -u +%Y-%m-%dT%T.%3NZ`] (sys) Starting" >> /var/log/gamblacoin.sys.log
	end script

pre-stop script
	rm /var/run/gamblacoin.pid
	echo "[`date -u +%Y-%m-%dT%T.%3NZ`] (sys) Stopping" >> /var/log/gamblacoin.sys.log
end script


